# scripts/server.service
[Unit]
Description=NIDS Central Server
After=network.target postgresql.service redis.service
Requires=postgresql.service redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/nids/server/server
EnvironmentFile=/etc/nids/server.env
ExecStart=/opt/nids/server/venv/bin/python app.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nids-server

[Install]
WantedBy=multi-user.target
